# deps/com.user.slsvc.plist
<?xml version="1.6" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.6">
<dict>
    <key>Label</key>
    <string>com.user.slsvc</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/slhealthupsvc</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/tmp/slsvc.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/slsvc.log</string>
</dict>
</plist>
